# scripts/launchagents/com.astromalik.cli.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.astromalik.cli.weekly</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ASTROMALIK_CLI_BIN</key>
        <string>/Users/eduardoariasbravo/Developer/AstroMalik-macOS/.build/release/astromalik-cli</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec "${ASTROMALIK_CLI_BIN:-/Users/eduardoariasbravo/Developer/AstroMalik-macOS/.build/release/astromalik-cli}" --chart "Edu" --scope weekly --model sonnet --output joplin:AstroMalik</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>7</integer>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/eduardoariasbravo/Library/Logs/AstroMalik/cli.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/eduardoariasbravo/Library/Logs/AstroMalik/cli.err</string>
</dict>
</plist>
